# com.CyberBrothers.akademitrack.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- This is what macOS shows in Background Items -->
  <key>Label</key>
  <string>AkademiTrack</string>
  
  <!-- Path to app bundle -->
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/open</string>
    <string>-a</string>
    <string>/Applications/AkademiTrack.app</string>
  </array>
  
  <!-- Run when the user logs in -->
  <key>RunAtLoad</key>
  <true/>
  
  <!-- Don't keep alive if it crashes -->
  <key>KeepAlive</key>
  <false/>
</dict>
</plist>